# vpn/walt/vpn/setup/walt-vpn-server.socket
[Socket]
ListenStream=/var/run/walt-vpn.sock

[Install]
WantedBy = sockets.target
