# vpn/walt/vpn/setup/walt-vpn-server.socket
[Socket]
ListenDatagram=/var/run/walt-vpn.sock
SocketUser=walt-vpn
SocketMode=0600

[Install]
WantedBy = sockets.target
