# server/walt/server/setup/walt-server-httpd.service
# walt-server-httpd - WalT platform HTTP daemon
#
# The WalT platform provides a lightweight distributed testbed for
# sensor networks, distributed protocols, distributed data
# management.

[Unit]
Description=WalT server HTTP service
# needs main daemon to be running
Requires=walt-server.service
Requires=walt-server-netconfig.service
After=walt-server-netconfig.service

[Service]
Environment=PYTHONUNBUFFERED=1
# After the initialization phase, we will notify systemd
# that we are started successfully.
Type=notify
NotifyAccess=all
ExecStart=walt-server-httpd
Restart=on-failure

[Install]
WantedBy=walt-server.service
