# server/walt/server/setup/walt-server-nbd.socket
[Socket]
ListenStream=10809

[Install]
WantedBy = sockets.target
